# monvpn.cron
@reboot $HOME/bin/monvpn >/var/tmp/monvpn.log 2>&1
